# overlay_qtota/usr/lib/systemd/system/qtota.service
[Unit]
Description=Qt OTA Client

[Service]
ExecStart=/usr/bin/qtota
EnvironmentFile=/etc/environment

[Install]
WantedBy=multi-user.target
